# res/wlr-gamma-service.service
[Unit]
Description=wlroots Brightness Control
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=wlr-gamma-service

[Install]
WantedBy=graphical-session.target
